# crontab
0 * * * * * root /root/image_syncer/image_sync.sh