# crontab
0-59 * * * *  root /bin/bash /root/image_syncer/image_sync.sh